# deploy/supervisor_almo.conf
[program:fit]
directory=/projects/almo/
command=/projects/almo_venv/bin/gunicorn app:app -b 127.0.0.1:7000
autostart=true
autorestart=true
